# deploy/com.dua.feature-forge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd plist to auto-start Feature Forge server on Mac Mini boot/login -->
<!-- Install: cp deploy/com.dua.feature-forge.plist ~/Library/LaunchAgents/ -->
<!--          launchctl load ~/Library/LaunchAgents/com.dua.feature-forge.plist -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dua.feature-forge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/doruntinaramadani/appmirror-tool-starter/server.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/doruntinaramadani/appmirror-tool-starter</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/feature-forge.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/feature-forge.error.log</string>
</dict>
</plist>
